# system/services/dhcpcd/run
#!/bin/sh
exec 2>&1

# build dhcpcd.conf
config_path=/run/dhcpcd.conf
for f in /etc/dhcpcd.conf.d/*; do
    cat "$f" >> $config_path
done

exec dhcpcd -f "$config_path" -B
