# runit-user/services/pipewire/run
#!/bin/sh
exec 2>&1

source ../../vars
export XDG_RUNTIME_DIR

exec pipewire
