# user/services/dbus/run
#!/bin/sh
exec 2>&1

source ../../vars

exec dbus-daemon --session --nofork --nopidfile --address "$DBUS_SESSION_BUS_ADDRESS"
