# runit-user/services/wireplumber/run
#!/bin/sh
exec 2>&1

source ../../lib/check-dependency
check_dependency pipewire

source ../../vars
export XDG_RUNTIME_DIR
export DBUS_SESSION_BUS_ADDRESS

exec wireplumber
