# system/services/rpcstatd/run
#!/bin/zsh
exec 2>&1

source ../../lib/check-dependency
check_dependency rpcbind

exec rpc.statd --foreground --no-syslog --no-notify
