# runit-system/services/sshd/run
#!/bin/sh
exec 2>&1

ssh-keygen -A > /dev/null
exec /usr/sbin/sshd -D
