# runit-system/services/mount-boot/run
#!/bin/sh
exec 2>&1

source ./config

mountpoint -q /boot || mount /dev/disk/by-partuuid/$partuuid /boot

# oneshot
exec rm -r $PWD
