# system/services/dbus/run
#!/bin/sh
exec 2>&1

mkdir -p "/run/dbus"
exec dbus-daemon --system --nofork --nopidfile
