# runit-system/services/zram/run
#!/bin/sh
exec 2>&1

set -e

source ./config

modprobe zram
zramctl /dev/zram0 --algorithm lzo-rle --size $size
mkswap -U clear /dev/zram0
swapon /dev/zram0
